# scripts/sentinel.service
[Unit]
Description=Sentinel System Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sentinel
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

[Install]
WantedBy=multi-user.target
